# dashboard.service
[Unit]
Description=Broadcast Device Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/axera/broadcast.axera
ExecStart=/usr/bin/env python3 /home/axera/broadcast.axera/dashboard.py
Restart=always
RestartSec=5
User=axera
Group=axera
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
